# unlocker-tool/app/src-tauri/LaunchDaemons/com.sofriendly.crosspoint.unlocker.helper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sofriendly.crosspoint.unlocker.helper</string>

    <key>BundleProgram</key>
    <string>Contents/MacOS/unlocker-helper</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>AssociatedBundleIdentifiers</key>
    <array>
        <string>com.sofriendly.crosspoint.unlocker</string>
    </array>

    <key>StandardOutPath</key>
    <string>/var/log/com.sofriendly.crosspoint.unlocker.helper.out.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/com.sofriendly.crosspoint.unlocker.helper.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>RUST_LOG</key>
        <string>info</string>
    </dict>
</dict>
</plist>
